# murkelhausen-data.service
[Unit]
Description=Prefect serves murkelhausen-data flow

[Service]
User=arkadius
WorkingDirectory=/home/arkadius/murkelhausen-data
ExecStart=/home/arkadius/murkelhausen-data/.venv/bin/python flow.py
Environment=PREFECT_LOGGING_LEVEL=DEBUG
Restart=always

[Install]
WantedBy=multi-user.target