# prefect.conf
[program:prefect]
autorestart = True
autostart = True
directory = /home/arkadius/
user = arkadius
command = /home/arkadius/.local/bin/prefect worker start --pool beowulf-local
